# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
logfile_maxbytes=10MB
loglevel=info
pidfile=/tmp/supervisord.pid

[unix_http_server]
file=/tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:llama-server]
command=/docker/run-llama.sh
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/llama.log
stderr_logfile=/var/log/llama.log
priority=10

[program:search-server]
command=node /app/search-server.mjs
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/search.log
stderr_logfile=/var/log/search.log
priority=15

[program:nginx]
# The sh wrapper polls the llama-server health endpoint before starting nginx.
# supervisord monitors the shell PID; if llama-server never becomes healthy
# the shell exits 1 and supervisord retries up to startretries times then stops.
# autorestart=unexpected means a clean nginx exit (e.g. SIGQUIT reload) is NOT restarted.
command=/bin/sh -c 'i=0; until curl -sf http://127.0.0.1:8080/health > /dev/null 2>&1; do sleep 2; i=$((i+1)); [ $i -ge 300 ] && echo "llama-server not ready after 10 min" && exit 1; done && /usr/sbin/nginx -g "daemon off;"'
autostart=true
autorestart=unexpected
startretries=3
stdout_logfile=/var/log/nginx-access.log
stderr_logfile=/var/log/nginx-error.log
priority=20
