# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
logfile_maxbytes=10MB
loglevel=info

[program:llama-server]
command=/docker/run-llama.sh
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/llama.log
stderr_logfile=/var/log/llama.log
priority=10

[program:nginx]
command=/bin/sh -c 'i=0; until curl -sf http://127.0.0.1:8080/health > /dev/null 2>&1; do sleep 2; i=$((i+1)); [ $i -ge 150 ] && echo "llama-server not ready after 5 min" && exit 1; done && /usr/sbin/nginx -g "daemon off;"'
autostart=true
autorestart=true
stdout_logfile=/var/log/nginx-access.log
stderr_logfile=/var/log/nginx-error.log
priority=20
