# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
logfile_maxbytes=10MB
loglevel=info

[program:llama-server]
command=/usr/local/bin/llama-server
    -m %(ENV_MODEL_PATH)s
    -c %(ENV_LLM_CONTEXT)s
    -ngl %(ENV_LLM_NGL)s
    -t %(ENV_LLM_THREADS)s
    --port 8080
    --host 127.0.0.1
    -fa
    --log-disable
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/llama.log
stderr_logfile=/var/log/llama.log
priority=10

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/var/log/nginx-access.log
stderr_logfile=/var/log/nginx-error.log
priority=20
